# cron.d/AATAMS_SATTAG_NRT
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au

* 0,12 * * * projectofficer AATAMS/AATAMS_sattag_nrt/aatams_sattag_nrt_main.sh
