# cron.d/AATAMS_SATTAG_DM
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au
0 0 * * 0,2,4,6 projectofficer . /mnt/ebs/data-services/env; AATAMS/AATAMS_sattag_dm/aatams_sattag_dm_sync.py
