# cron.d/AATAMS_SATTAG_DM
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au

0 0,12 * * * projectofficer AATAMS/AATAMS_sattag_dm/aatams_sattag_dm_main.sh
