# cron.d/ACORN
SHELL=/bin/bash

MAILTO=dan.fruehauf@utas.edu.au

# build hierarchy for incoming ACORN files
*/5 * * * * dfruehauf ACORN/build-dir-hierarchy/build-dir-hierarchy.sh /mnt/imos-t4/IMOS/staging/ACORN/acorn-migration /mnt/imos-t4/IMOS/staging/ACORN/acorn-migration-hierarchy
