# cron.d/ACORN
SHELL=/bin/bash

MAILTO=dan.fruehauf@utas.edu.au

# build hierarchy for incoming ACORN files
*/5 * * * * ggalibert ACORN/build-dir-hierarchy/build-dir-hierarchy.sh $INCOMING_DIR/ACORN/acorn-migration $INCOMING_DIR/ACORN/acorn-migration-hierarchy
