# cron.d/ARGO
MAILTO=Angus.Scheibner@utas.edu.au

0 18 * * 5 projectofficer ARGO/rsync.sh
