# cron.d/ARGO
MAILTO=Laurent.Besnard@utas.edu.au

0 20 * * * projectofficer . /mnt/ebs/data-services/env; ARGO/rsync.sh
