# cron.d/ARGO
MAILTO=Angus.Scheibner@utas.edu.au

0 20 * * * projectofficer ARGO/rsync.sh
