# cron.d/ARGO
MAILTO=craig.jones@utas.edu.au

0 18 * * 5 nobody ARGO/lftp_sync.sh
