# cron.d/FAIMMS_RSS
MAILTO=laurent.besnard@utas.edu.au

0 22 * * * lbesnard FAIMMS/FAIMMS_data_rss_channels_process_matlab/FAIMMS.sh
