# cron.d/FAIMMS_RSS
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au

0 5,11,17,23  * * * projectofficer FAIMMS/FAIMMS_data_rss_channels_process/FAIMMS.sh
