# cron.d/FAIMMS_RSS
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au

0 * * * * projectofficer FAIMMS/REALTIME/faimms.py
