# cron.d/FAIMMS_RSS
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au

0 0,8,16 * * * projectofficer . /mnt/ebs/data-services/env; FAIMMS/REALTIME/faimms.py
