# cron.d/FAIMMS_RSS
MAILTO=laurent.besnard@utas.edu.au

0 22 * * * lbesnard FAIMMS/FAIMMS_data_rss_channels_process_matlab/FAIMMS.sh && rsync --size-only --itemize-changes --stats -uhvrD  --progress /mnt/imos-t4/project_officers/wip/FAIMMS/data_opendap_folder_rsync/opendap/FAIMMS/  /mnt/imos-t3/IMOS/opendap/FAIMMS/
