# cron.d/LINK_CHECKER
SHELL=/bin/bash

5 8 * * 6 projectofficer bin/geonetwork_broken_link_checker.py --geonetwork-url  http://catalogue-imos.aodn.org.au/geonetwork --email-to natalia.atkins@utas.edu.au --filter-tag gmd:distributionInfo

