# cron.d/NRMN
SHELL=/bin/bash

MAILTO=benedicte.pasquer@utas.edu.au

0 21 * * 5 projectofficer NRMN/extract.sh
