# cron.d/NRS_DAR_YON_RSS
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au

0 * * * * projectofficer  ANMN/NRS_AIMS/REALTIME/NRS.sh
