# cron.d/NRS_MAI_NSI_RT
# Process real-time data from Maria Island and North Stradbroke Island
# NRS stations

SHELL=/bin/bash
MAILTO=marty.hidas@utas.edu.au

# m h  dom mon dow user    command
42 */6 * * *         nobody  ANMN/NRSrealtime/NRSrealtimeUpdate.sh
