# cron.d/NRS_DAR_YON_RSS
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au

0 0,6,12,18 * * * projectofficer  ANMN/NRS_AIMS_Darwin_Yongala_data_rss_channels_process/NRS.sh
