# cron.d/NRS_DAR_YON_RSS
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au

0 7,13,18 * * * projectofficer . /mnt/ebs/data-services/env; ANMN/NRS_AIMS/REALTIME/anmn_nrs_aims.py
