# cron.d/SOOP_TRV_RSS
MAILTO=laurent.besnard@utas.edu.au
0 20 * * * lbesnard SOOP/SOOP_TRV_aims_download_process_matlab/SOOP_TRV.sh
