# cron.d/SOOP_XBT
SHELL=/bin/bash
MAILTO=laurent.besnard@utas.edu.au
# m h  dom mon dow   command

# start XBT data
0 5,10,15,20  * * * lbesnard SOOP/SOOP_SST_ASF_XBT/main.sh XBT

