# cron.d/SOOP_BA
SHELL=/bin/bash
# PROCESSED data moved every week from STAGING to OPENDAP and PUBLIC
# m h  dom mon dow   command

MAILTO=benedicte.pasquer@utas.edu.au

# MOVE DATA daily FROM STAGING TO EITHER OPENDAP OR PUBLIC
15 1 * * *  bpasquer SOOP/SOOP_BA_filecopy_bash/BASH/moveSOOP_BA.sh

# MOVE RAW DATA WEEKLY FROM STAGING TO ARCHIVE
#15 1 * * 0  bpasquer SOOP/SOOP_BA_filecopy_bash/BASH/moveSOOP_BA_raw.sh
