# cron.d/SOOP_BA
SHELL=/bin/bash
# PROCESSED data moved every week from STAGING to OPENDAP and PUBLIC
# m h  dom mon dow   command  

MAILTO=benedicte.pasquer@utas.edu.au
#MOVE PROCESSED DATA daily FROM STAGING TO OPENDAP AND PUBLIC
15 1 * * *  bpasquer SOOP_BA/BASH/moveSOOP_BA_PROCESSED.sh

#MOVE RAW DATA daily FROM STAGING TO ARCHIVE
30 1 * * * bpasquer SOOP_BA/BASH/moveSOOP_BA_RAW.sh

