# cron.d/SOOP_ASF_SST
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au

0 12  * * * projectofficer . /mnt/ebs/data-services/env; SOOP/SOOP_ASF_SST_XBT/SOOP_BOM_ASF_SST.py
