# cron.d/SOOP_BA
SHELL=/bin/bash
# PROCESSED data moved every week from STAGING to OPENDAP and PUBLIC
# m h  dom mon dow   command

MAILTO=benedicte.pasquer@utas.edu.au

# MOVE DATA daily FROM STAGING TO EITHER OPENDAP OR PUBLIC
15 1 * * *  bpasquer SOOP_BA/BASH/moveSOOP_BA.sh
