# cron.d/SRS_OC_LJCO_AERONET
MAILTO=laurent.besnard@utas.edu.au

0 22 * * * lbesnard pyton SRS/SRS_OC_LJCO_AERONET/downloadAeronetData.py