# cron.d/SRS_OC_LJCO_AERONET
SHELL=/bin/bash

MAILTO=laurent.besnard@utas.edu.au

0 19 * * * projectofficer SRS/SRS_OC_LJCO_AERONET/main.sh
